# resources/waybar.service.in
[Unit]
Description=Highly customizable Wayland bar for Sway and Wlroots based compositors.
Documentation=https://github.com/Alexays/Waybar/wiki/
PartOf=graphical-session.target
After=graphical-session.target

[Service]
ExecStart=@prefix@/bin/waybar
Restart=on-failure

[Install]
WantedBy=graphical-session.target
